# apps/.config/systemd/user/app-com.mitchellh.ghostty.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/ghostty --gtk-single-instance=false
